# logstash/templates/etc/kibana/kibanainit.rb.erb
#!/bin/bash
#
#### BEGIN INIT INFO
## Provides:          kibana
## Required-Start:    $elasticserach
## Required-Stop:
## Default-Start:     2 3 4 5
## Default-Stop:      0 1 6
## Short-Description: Starts Kibana
## Description:       Starts Kibana using start-stop-daemon
#### END INIT INFO
#
# Author - Matt Mencel - 2012
# Borrowed heavily from Patrick McKenzie's example
# http://www.kalzumeus.com/2010/01/15/deploying-sinatra-on-ubuntu-in-which-i-employ-a-secretary/
#
# kibana-daemon       Startup script for the Kibana Web UI running on sinatra
#
## Copy to /etc/init.d/kibana and make it executable
## Add to system startup through chkconfig (CentOS/RHEL) or Upstart (Ubuntu)
KIBANA_PATH="/opt/kibana"
cd $KIBANA_PATH
ruby $KIBANA_PATH/kibana-daemon.rb $1
RETVAL=$?

exit $RETVAL